# roles/prometheus_node_exporter/templates/init.service.j2
[Unit]
Description={{serviceName}}
Wants=network-online.target
After=network-online.target

[Service]
User={{ userId }}
Group={{ groupId }}
Restart=on-failure
Type=simple
ExecStart={{ exec_command }}

[Install]
WantedBy=multi-user.target